# 95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/julienvallee/bin/monitor-hotplug.sh"
